# extracted_firmwares/squashfs-root/usr/lib/opkg/info/php7-fastcgi.list
/etc/init.d/php7-fastcgi
/etc/config/php7-fastcgi
